# config/crontab
* * * * * /bin/bash -l -c 'cd /planeray && RAILS_ENV=production bundle exec rake events:dumps --silent'
